# deploy/vouchernet.service
[Unit]
Description=VoucherNet API + frontend
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=vouchernet
Group=vouchernet
WorkingDirectory=/var/www/vouchernet/artifacts/api-server
EnvironmentFile=/var/www/vouchernet/.env
Environment=NODE_ENV=production
Environment=NODE_OPTIONS=--max-old-space-size=2048
ExecStart=/usr/bin/node dist/index.js
Restart=always
RestartSec=5
TimeoutStopSec=15
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vouchernet

# Sécurité minimale
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
